# etherrat_behavior.yml
# Sigma — Distributed EtherHiding Resolver Kit (EtherRAT family)
# Report: CTI-2026-0614-ETHERHIDING-KIT  |  TLP:CLEAR
# Multiple rules in one file (----separated YAML documents)

title: Conhost Headless Spawning Node or Cmd From User Path (EtherRAT chain)
id: a1f3e2c4-7b6d-4e9a-9c10-etherrat0001
status: experimental
description: >
  Detects conhost.exe launched with --headless executing node.exe or cmd.exe,
  the indirect-execution pattern used by the EtherHiding/EtherRAT loader to run
  its Node.js C2 resolver stage while hiding the console window.
references:
  - CTI-2026-0614-ETHERHIDING-KIT
author: Independent CTI
date: 2026-06-14
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\conhost.exe'
    CommandLine|contains: '--headless'
  child:
    CommandLine|contains:
      - '\node.exe'
      - '\cmd'
  condition: selection and child
falsepositives:
  - Legitimate developer tooling using conhost headless mode (rare)
level: high
tags:
  - attack.defense_evasion
  - attack.t1202
  - attack.execution
---
title: Node.exe Executing XML-Disguised Script Argument (EtherRAT resolver)
id: b2c4d6e8-1a3f-4c5b-8d9e-etherrat0002
status: experimental
description: >
  Detects node.exe invoked with a .xml file as its script argument from a
  user-writable AppData path - the EtherHiding resolver stage masquerade.
references:
  - CTI-2026-0614-ETHERHIDING-KIT
author: Independent CTI
date: 2026-06-14
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\node.exe'
    CommandLine|contains: '.xml'
  path:
    CommandLine|contains:
      - '\AppData\Local\'
      - '\AppData\Roaming\'
  condition: selection and path
falsepositives:
  - Unusual; node rarely takes a .xml as a script argument
level: high
tags:
  - attack.execution
  - attack.t1059.007
  - attack.t1036.005
---
title: Node.exe Masqueraded in Oracle Java Path (EtherRAT)
id: c3d5e7f9-2b4a-4d6c-9e0f-etherrat0003
status: experimental
description: >
  Detects node.exe located under the Oracle\Java\javapath directory, a
  masquerade location used by the EtherHiding kit loader.
references:
  - CTI-2026-0614-ETHERHIDING-KIT
author: Independent CTI
date: 2026-06-14
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|contains|all:
      - '\Oracle\Java\javapath\'
      - 'node.exe'
  condition: selection
falsepositives:
  - None expected; node.exe does not legitimately install to the Java path
level: high
tags:
  - attack.defense_evasion
  - attack.t1036.005
---
title: MSI-Spawned Run Key Persistence To Random AppData Binary (EtherRAT)
id: d4e6f8a0-3c5b-4e7d-af10-etherrat0004
status: experimental
description: >
  Detects reg.exe writing a Run-key autostart entry pointing into a randomly
  named %AppData%\Local subdirectory, the persistence step of the kit loader.
references:
  - CTI-2026-0614-ETHERHIDING-KIT
author: Independent CTI
date: 2026-06-14
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\reg.exe'
    CommandLine|contains|all:
      - '\CurrentVersion\Run'
      - '\AppData\Local\'
  condition: selection
falsepositives:
  - Some legitimate installers; correlate with conhost/node activity
level: medium
tags:
  - attack.persistence
  - attack.t1547.001
